# shutdowntimer.service
[Unit]
Description= apagado a hora específica
After=network-online.target


[Service]
Type=simple
ExecStart=shutdown -P now
